# files/etc/uci-defaults/luci-dashboard
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@dashboard[-1]
	add ucitrack dashboard
	set ucitrack.@dashboard[-1].init=dashboard
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
